# config/go2rtc.service
[Unit]
Description=go2rtc — RTSP/WebRTC bridge for rover cameras
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gimli1mb
WorkingDirectory=/opt/gimli
ExecStart=/usr/local/bin/go2rtc -c /opt/gimli/config/go2rtc.yaml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
